# systemd/vivobook-sound-fix.service
[Unit]
Description=ASUS Vivobook S14 Flip TP3407SA speaker fix
After=sound.target

[Service]
Type=oneshot
ExecStartPre=/bin/bash -c 'until [ -e /dev/snd/controlC0 ]; do sleep 0.5; done'
ExecStart=/usr/bin/asus-sound-fix
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
